# templates/default/airflow-scheduler.service.erb
# Copyright 2015 Sergey Bahchissaraitsev

# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at

#     http//www.apache.org/licenses/LICENSE-2.0

# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

[Unit]
Description=Airflow scheduler daemon
After=network.target <%= @deps %>

[Service]
ExecStart = /usr/bin/docker run --rm --name airfow-scheduler \
           -v <%= node['airflow']['base_dir'] %>/airflow.cfg:/airflow/airflow.cfg \
           -v <%= node['airflow']['data_volume']['dags_dir'] %>:/airflow/dags \
           -v <%= node['airflow']['data_volume']['log_dir'] %>:/airflow/logs \
           -v <%= node['airflow']['data_volume']['secrets_dir'] %>:/airflow/secrets \
           --network=host \
           <%= @image_name %> \
           airflow scheduler
ExecStop = /usr/bin/docker stop airflow-scheduler

Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
